# deploy/pixelrag-agent.service
[Unit]
Description=PixelRAG Agent backend (Claude Agent SDK, subscription auth)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yichuan
Group=yichuan
WorkingDirectory=/home/yichuan/visrag/web
Environment=HOME=/home/yichuan
Environment=AGENT_PORT=30010
Environment=PIXELRAG_SEARCH_URL=http://localhost:30001
Environment=CHAT_MAX_BUDGET_USD=2.00
Environment=ALLOWED_ORIGIN=https://pixelrag.ai
Environment=RL_PER_IP=12
Environment=RL_GLOBAL_DAILY=500
Environment=RL_MAX_CONCURRENT=3
ExecStart=/home/yichuan/.nix-profile/bin/node /home/yichuan/visrag/web/agent-server.mjs
Restart=always
RestartSec=3
StandardOutput=append:/home/yichuan/visrag/logs/agent-server.log
StandardError=append:/home/yichuan/visrag/logs/agent-server.log

[Install]
WantedBy=multi-user.target
